# com.openclaw.agent-console.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.agent-console</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>[workspace-root]/start-agent-console.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>[workspace-root]</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>[workspace-root]/console.log</string>

  <key>StandardErrorPath</key>
  <string>[workspace-root]/console.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>~</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOST</key>
    <string>0.0.0.0</string>
    <key>PORT</key>
    <string>8787</string>
  </dict>
</dict>
</plist>